# fleetctl/templates/docker_cleanup.timer
[Unit]
Description=Weekly cleanup of stale docker images
Requires=docker.service

[Timer]
OnCalendar=weekly

[X-Fleet]
Global=true
